# deploy/dchat-socket.service
[Unit]
Description=Dchat Socket.io Server
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/dchat/backend
Environment="PATH=/home/ec2-user/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ec2-user/dchat/backend"
ExecStart=/usr/bin/python3 -m src.socket_app
Restart=always
RestartSec=10
StandardOutput=append:/var/log/dchat/socket.log
StandardError=append:/var/log/dchat/socket-error.log

[Install]
WantedBy=multi-user.target
